# udev-rules/99-azure-network-tuned.rules
# Azure Network Adapter UDEV Rules
# Optimizes performance for Azure network interfaces on VM startup
#
# Installation:
#   sudo cp 99-azure-network-tuned.rules /etc/udev/rules.d/
#   sudo chmod 644 /etc/udev/rules.d/99-azure-network-tuned.rules
#
# Apply without reboot:
#   sudo udevadm control --reload-rules && udevadm trigger
#
# Performance tuning notes:
# Ring buffer sizes (RX/TX): 1024 is balanced starting point
# Mellanox documentation recommends 8192/8192 for high throughput
# Optimal values depend on workload - test for your environment

# Accelerated Networking (MANA/Mellanox) - Set ring buffers to 1024 
# Applies to MANA (hv_pci) accelerated interfaces
SUBSYSTEM=="net", DRIVERS=="hv_pci", ACTION=="add", RUN+="/usr/sbin/ethtool -G $env{INTERFACE} rx 1024 tx 1024"

# Synthetic Networking (hv_netvsc) - Set ring buffers to 1024
# Applies to standard synthetic interfaces
SUBSYSTEM=="net", DRIVERS=="hv_netvsc*", ACTION=="add", RUN+="/usr/sbin/ethtool -G $env{INTERFACE} rx 1024 tx 1024"

# Synthetic Networking - Disable UDP hash on synthetic interface
# Prevents UDP traffic hashing which can cause packet reordering
SUBSYSTEM=="net", DRIVERS=="hv_netvsc*", ACTION=="add", RUN+="/usr/sbin/ethtool -N $env{INTERFACE} rx-flow-hash udp4 sd"

# Accelerated Networking - Use NOQUEUE qdisc for lower latency
# NOQUEUE bypasses the queuing discipline for faster packet processing
SUBSYSTEM=="net", DRIVERS=="hv_pci|mlx*", ACTION=="add", RUN+="/usr/sbin/tc qdisc replace dev $env{INTERFACE} root noqueue"

# Set transmit queue length for all network interfaces
# Higher values: better for bursty traffic
# Lower values: lower latency, reduced memory overhead
SUBSYSTEM=="net", DRIVERS=="hv_netvsc*|mlx*|hv_pci", ACTION=="add", ATTR{tx_queue_len}="2048"
SUBSYSTEM=="net", KERNEL=="lo", ACTION=="add", ATTR{tx_queue_len}="5000"
